# service/voice-to-text.service
[Unit]
Description=Voice-to-Text D-Bus Service
After=graphical-session.target pipewire.service wireplumber.service
Wants=graphical-session.target pipewire.service wireplumber.service

[Service]
Type=exec
ExecStart=%h/.local/bin/voice-to-text-dbus
Restart=on-failure
RestartSec=3
# Import session environment so clipboard tools (wl-copy/xclip) can find
# WAYLAND_DISPLAY / DISPLAY and dbus finds the session bus
PassEnvironment=DISPLAY WAYLAND_DISPLAY XDG_SESSION_TYPE VOXTRAL_API_KEY DEEPGRAM_API_KEY GROQ_API_KEY

[Install]
WantedBy=default.target
